# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:wa-bridge]
command=node /app/wa-bridge/src/index.js
directory=/app/wa-bridge
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOST="0.0.0.0",PORT="3001",BACKEND_URL="http://localhost:8000",PUPPETEER_SKIP_CHROMIUM_DOWNLOAD="true",PUPPETEER_EXECUTABLE_PATH="/usr/bin/chromium"

[program:backend]
command=/bin/bash -c "/app/venv/bin/alembic upgrade head && exec /app/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000"
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DB_URL="sqlite:////app/data/wishing_bot.db",WA_BRIDGE_URL="http://localhost:3001",FRONTEND_ORIGIN="http://%(ENV_NAS_IP)s:8080"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
